# systemd/cupcake-backend.service
[Unit]
Description=Cupcake Django Backend (Gunicorn)
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=cupcake
Group=cupcake
WorkingDirectory=/opt/cupcake/backend
EnvironmentFile=/opt/cupcake/.env
ExecStart=/opt/cupcake/venv/bin/gunicorn cupcake_vanilla.asgi:application \
    --bind 127.0.0.1:8000 \
    --workers 4 \
    --timeout 300 \
    -k uvicorn.workers.UvicornWorker \
    --log-level info \
    --access-logfile /var/log/cupcake/gunicorn-access.log \
    --error-logfile /var/log/cupcake/gunicorn-error.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
